# conf/cloud9.conf
[program:cloud9]
command = node /c9sdk/server.js -w /workspace --listen 0.0.0.0
redirect_stderr=true
environment = NODE_ENV="production"
