# conf/nodervisor.conf
[program:nodervisor]
command = node /nodervisor/app.js
redirect_stderr=true
