# conf/vnc.conf
[program:vnc]
command = bash -c "vncserver :1 -geometry 1280x800 -depth 24 && tail -F /root/.vnc/*.log"
redirect_stderr=true
