# conf/wrapdocker.conf
[program:wrapdocker]
command=/usr/local/bin/wrapdocker
redirect_stderr=true
